# stage6/10-install-i2c/files/openplotter-i2c-read.service
[Service]
ExecStart=openplotter-i2c-read
StandardOutput=syslog
StandardError=syslog
